# {{cookiecutter.ansible_project_name}}/roles/{{cookiecutter.vm_slug}}/templates/webapp/celery_supervisor.conf
[program:{{ app_name }}_celery]
command = {{ app_virtualenv_dir }}/bin/celery_start
user = {{ app_user }}
stdout_logfile = {{ app_log_dir }}/celery.log
stderr_logfile = {{ app_log_dir }}/celery.log
