# {{cookiecutter.ansible_project_name}}/roles/{{cookiecutter.vm_slug}}/templates/webapp/gunicorn_supervisor.conf
[program:{{ app_name }}]
command = {{ app_virtualenv_dir }}/bin/gunicorn_start
user = {{ app_user }}
stdout_logfile = {{ app_log_dir }}/gunicorn.log
stderr_logfile = {{ app_log_dir }}/gunicorn.log
